# ansible/roles/sharedfolder/files/mpi.mount
[Unit]
Description=mount mpi

[Mount]
What=pi0:/mpi
Where=/mpi
Type=nfs