# ansible/roles/sharedfolder/files/mpi.automount
[Unit]
Description=automount mpi

[Automount]
Where=/mpi

[Install]
WantedBy=multi-user.target